# supervisord.conf
[program:pratter]
command = /home/dotcloud/current/script/dotcloud.sh
